# devterm-a06/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/devterm-printer
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/devterm-audio-patch
-rwxr-xr-x	root/root	etc/rc.d/devterm-fan-daemon
-rwxr-xr-x	root/root	etc/rc.d/devterm-printer
-rwxr-xr-x	root/root	etc/rc.d/devterm-socat
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/DevTerm_keyboard_firmware_v0.1_utils.sh
-rwxr-xr-x	root/root	usr/bin/DevTerm_keyboard_firmware_v0.21_utils.sh
-rwxr-xr-x	root/root	usr/bin/DevTerm_keyboard_firmware_v0.2_utils.sh
-rwxr-xr-x	root/root	usr/bin/DevTerm_keyboard_firmware_v0.2a_utils.sh
-rwxr-xr-x	root/root	usr/bin/DevTerm_keyboard_firmware_v0.3_utils.sh
-rwxr-xr-x	root/root	usr/bin/audio_3.5_patch.py
-rwxr-xr-x	root/root	usr/bin/devterm-a06-gearbox
-rwxr-xr-x	root/root	usr/bin/devterm_socat_perms
-rwxr-xr-x	root/root	usr/bin/devterm_thermal_printer.elf
-rwsr-xr-x	root/root	usr/bin/gpio
-rwxr-xr-x	root/root	usr/bin/temp_fan_daemon_a06.py
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/ads1115.h
-rw-r--r--	root/root	usr/include/bmp180.h
-rw-r--r--	root/root	usr/include/drcNet.h
-rw-r--r--	root/root	usr/include/drcSerial.h
-rw-r--r--	root/root	usr/include/ds1302.h
-rw-r--r--	root/root	usr/include/ds18b20.h
-rw-r--r--	root/root	usr/include/gertboard.h
-rw-r--r--	root/root	usr/include/htu21d.h
-rw-r--r--	root/root	usr/include/lcd.h
-rw-r--r--	root/root	usr/include/lcd128x64.h
-rw-r--r--	root/root	usr/include/max31855.h
-rw-r--r--	root/root	usr/include/max5322.h
-rw-r--r--	root/root	usr/include/maxdetect.h
-rw-r--r--	root/root	usr/include/mcp23008.h
-rw-r--r--	root/root	usr/include/mcp23016.h
-rw-r--r--	root/root	usr/include/mcp23016reg.h
-rw-r--r--	root/root	usr/include/mcp23017.h
-rw-r--r--	root/root	usr/include/mcp23s08.h
-rw-r--r--	root/root	usr/include/mcp23s17.h
-rw-r--r--	root/root	usr/include/mcp23x08.h
-rw-r--r--	root/root	usr/include/mcp23x0817.h
-rw-r--r--	root/root	usr/include/mcp3002.h
-rw-r--r--	root/root	usr/include/mcp3004.h
-rw-r--r--	root/root	usr/include/mcp3422.h
-rw-r--r--	root/root	usr/include/mcp4802.h
-rw-r--r--	root/root	usr/include/pcf8574.h
-rw-r--r--	root/root	usr/include/pcf8591.h
-rw-r--r--	root/root	usr/include/piFace.h
-rw-r--r--	root/root	usr/include/piGlow.h
-rw-r--r--	root/root	usr/include/piNes.h
-rw-r--r--	root/root	usr/include/pseudoPins.h
-rw-r--r--	root/root	usr/include/rht03.h
-rw-r--r--	root/root	usr/include/scrollPhat.h
-rw-r--r--	root/root	usr/include/sn3218.h
-rw-r--r--	root/root	usr/include/softPwm.h
-rw-r--r--	root/root	usr/include/softServo.h
-rw-r--r--	root/root	usr/include/softTone.h
-rw-r--r--	root/root	usr/include/sr595.h
-rw-r--r--	root/root	usr/include/wiringCPi.h
-rw-r--r--	root/root	usr/include/wiringPi.h
-rw-r--r--	root/root	usr/include/wiringPiI2C.h
-rw-r--r--	root/root	usr/include/wiringPiSPI.h
-rw-r--r--	root/root	usr/include/wiringSerial.h
-rw-r--r--	root/root	usr/include/wiringShift.h
-rw-r--r--	root/root	usr/include/wpiExtensions.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/filter/
-rwxr-xr-x	root/root	usr/lib/cups/filter/rastertocpi
lrwxrwxrwx	root/root	usr/lib/libwiringPi.so -> libwiringPi.so.2.50
-rwxr-xr-x	root/root	usr/lib/libwiringPi.so.2.50
lrwxrwxrwx	root/root	usr/lib/libwiringPiDev.so -> libwiringPiDev.so.2.50
-rwxr-xr-x	root/root	usr/lib/libwiringPiDev.so.2.50
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/wiringpid
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/cups/
drwxr-xr-x	root/root	usr/share/cups/model/
drwxr-xr-x	root/root	usr/share/cups/model/clockworkpi/
-rw-r--r--	root/root	usr/share/cups/model/clockworkpi/cpi58.ppd
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/gpio.1.gz
